# cellular-keepalive.service
[Unit]
Description=Cellular Keepalive Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/sbin/cellular-keepalive
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

